# build/mchn/PT/pidp10.service
[Unit]
Description=PiDP-10
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/pidp10
ExecStart=/bin/sh ./pidp10.sh start
ExecStop=/bin/sh ./pidp10.sh stop
Restart=always
RestartSec=10
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
